# autostart/xmodmap.desktop
[Desktop Entry]
Exec=xmodmap ~/.xmodmap
